# scripts/gunicorn_mimic_upstart.conf
description "mimic"

start on (filesystem)
stop on runlevel [016]

respawn
console log
setuid nobody
setgid nogroup
chdir <mimic folder that has wsgi.py file>

exec gunicorn wsgi:application